# root/etc/services.d/dbus/run
#!/usr/bin/with-contenv bash

# Start dbus
dbus-daemon --config-file=/usr/share/dbus-1/system.conf --print-address --nofork

